# www.luokr.com/V0.3/doc/supervisord.conf.d/www.luokr.com.conf
[program:www.luokr.com] 
directory=/var/www/www.luokr.com/V0.3/app 
command=/bin/env python svr.py
autostart=true
